# config/upnp-map.timer
[Unit]
Description=Renew UPnP mapping every hour

[Timer]
OnBootSec=30
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
