# config/upnp-map.service
[Unit]
Description=UPnP Port Mapping
After=network-online.target xray.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/upnp-map.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
